# templates/default/bitcoind.upstart.erb
# Bitcoind Upstart script

description "<%= @bitcoind_name %>"

start on filesystem
stop on runlevel [016]

setuid <%= @bitcoind_user %>
setgid <%= @bitcoind_group %>

env HOME=<%= @bitcoind_home %>

exec <%= File.join(@bitcoind_bin_location, @bitcoind_bin_name) %> -datadir=<%= @bitcoind_data_dir %> <%= @bitcoind_options %> 2>&1 | logger -t <%= @bitcoind_bin_name %>
