# templates/default/bitcoind.upstart.erb
# Bitcoind Upstart script

description "Bitcoind"

start on filesystem
stop on runlevel [016]

setuid <%= @bitcoind_user %>
setgid <%= @bitcoind_group %>

env HOME=<%= @bitcoind_home %>

exec <%= @bitcoind_bin_location %><%= @bitcoind_bin_name %> <%= @bitcoind_options %> 2>&1 | logger -t bitcoind
